# spring_ksm/util_spring/sts-bundle/pivotal-tc-server/developer-4.1.0.RELEASE/templates/systemd/conf/tcruntime.service
[Unit]
Description=tc Runtime Instance ${instance.name}
After=network.target

[Service]
User=${runtime.user:tcserver}
Type=forking
WorkingDirectory=${installation.directory}
ExecStart=${installation.directory}/tcserver start ${instance.name} -i ${instances.directory}
TimeoutStartSec=0

[Install]
WantedBy=default.target

